# config/ubuntu-12.04/cron.daily/guidcoin-prod-database
#!/bin/bash

echo 'VACUUM ANALYZE;' | psql guidcoin_prod
/var/www/python/guidcoin-prod/venv/bin/python /var/www/python/guidcoin-prod/web/manage_prod.py clearsessions

